# ansible/roles/deploy/templates/backbeard-supervisor.conf
[program:backbeard]
command=/tmp/backbeard-staging/bin/backbeard -Dplay.http.secret.key=abcdefghijk
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=true
